# scripts/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/tmp/django.sock

[Install]
WantedBy=sockets.target